# src/healthchecks_cleanup.cron
0 0 1 */3 * /code/healthchecks_cleanup.sh
